# comfy-tunnel.service
[Unit]
Description=SSH Tunnel to Digital Ocean for ComfyUI Media App
After=network.target

[Service]
Type=simple
User=Me
ExecStart=/usr/bin/ssh -i /c/Users/Me/.ssh/comfy_tunnel -R 8000:localhost:8000 -N -o ServerAliveInterval=60 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes root@192.81.208.72
Restart=always
RestartSec=10
Environment="AUTOSSH_GATETIME=0"

[Install]
WantedBy=multi-user.target